# src/pmproxy/pmproxy.service.in
[Unit]
Description=Proxy for Performance Metrics Collector Daemon
Documentation=man:pmproxy(1)
After=network-online.target pmcd.service redis.service avahi-daemon.service

[Service]
Type=notify
NotifyAccess=all
Restart=always
ExecStart=@PCP_RC_DIR@/pmproxy start-systemd
Group=@PCP_GROUP@
User=@PCP_USER@

[Install]
WantedBy=multi-user.target
