# src/pmproxy/pmproxy.service.in
[Unit]
Description=Proxy for Performance Metrics Collector Daemon
Documentation=man:pmproxy(1)
After=network-online.target valkey.service redis.service avahi-daemon.service pmcd.service

[Service]
Type=notify
NotifyAccess=all
Restart=always
ExecStart=@PCP_SYSCONF_DIR@/pmproxy/rc start-systemd
Group=@PCP_GROUP@
User=@PCP_USER@

[Install]
WantedBy=multi-user.target
