# src/pmproxy/pmproxy.service.in
[Unit]
Description=Proxy for Performance Metrics Collector Daemon
Documentation=man:pmproxy(8)
After=local-fs.target network.target pmcd.service

[Service]
Type=oneshot
ExecStart=@path@/pmproxy start
ExecStop=@path@/pmproxy stop
RemainAfterExit=yes 

[Install]
WantedBy=multi-user.target
