# src/pmmgr/pmmgr.service.in
[Unit]
Description=Performance Metrics Daemon Manager
Documentation=man:pmmgr(8)
Wants=avahi-daemon.service
After=network.target avahi-daemon.service pmcd.service

[Service]
Type=forking
Restart=always
ExecStart=@PCP_RC_DIR@/pmmgr start
ExecStop=@PCP_RC_DIR@/pmmgr stop
PIDFile=@PCP_RUN_DIR@/pmmgr.pid

[Install]
WantedBy=multi-user.target
