# src/pmie/pmie.service.in
[Unit]
Description=Performance Metrics Inference Engine
Documentation=man:pmie(1)
After=network-online.target pmcd.service
Before=pmie_check.timer pmie_daily.timer
BindsTo=pmie_check.timer pmie_daily.timer
Wants=pmcd.service

[Service]
Type=notify
Restart=always
NotifyAccess=all
TimeoutSec=120

ExecStart=@PCP_RC_DIR@/pmie start-systemd
ExecStop=@PCP_RC_DIR@/pmie stop-systemd
PIDFile=@PCP_RUN_DIR@/pmie.pid

[Install]
WantedBy=multi-user.target
