# src/pmie/pmie_farm_check.service.in
[Unit]
Description=Check and migrate non-primary pmie farm instances
Documentation=man:pmiectl(1)
ConditionPathExists=!@CRONTAB_PATH@

[Service]
Type=@SD_SERVICE_TYPE@
Restart=no
TimeoutStartSec=4h
TimeoutStopSec=120
ExecStart=@PCP_BIN_DIR@/pmiectl -m check
WorkingDirectory=@PCP_VAR_DIR@
# root so pmiectl can migrate pmie processes to the pmie_farm service
Group=root
User=root
