# src/pmie/pmie_daily.service.in
[Unit]
Description=Process PMIE logs
Documentation=man:pmie_daily(1)
ConditionPathExists=!@CRONTAB_PATH@

[Service]
Type=@SD_SERVICE_TYPE@
Restart=no
TimeoutStartSec=1h
Environment="PMIE_DAILY_PARAMS=-X xz -x 3"
EnvironmentFile=-@PCP_SYSCONFIG_DIR@/pmie_timers
ExecStart=@PCP_BINADM_DIR@/pmie_daily $PMIE_DAILY_PARAMS
WorkingDirectory=@PCP_VAR_DIR@
Group=@PCP_GROUP@
User=@PCP_USER@
