# src/pmie/pmie.service.in
[Unit]
Description=Performance Metrics Inference Engine
Documentation=man:pmie(1)
After=network-online.target pmcd.service
After=pmie_check.timer pmie_daily.timer
BindsTo=pmie_check.timer pmie_daily.timer

[Service]
Type=forking
Restart=always
ExecStart=@PCP_RC_DIR@/pmie start
ExecStop=@PCP_RC_DIR@/pmie stop
PIDFile=@PCP_RUN_DIR@/pmie.pid

[Install]
WantedBy=multi-user.target
