# src/pmie/pmie_check.service.in
[Unit]
Description=Check PMIE instances are running
Documentation=man:pmie_check(1)
ConditionPathExists=!@CRONTAB_PATH@

[Service]
Type=oneshot
Environment="PMIE_CHECK_PARAMS=-C"
EnvironmentFile=-@PCP_SYSCONFIG_DIR@/pmie_timers
ExecStart=@PCP_BINADM_DIR@/pmie_check $PMIE_CHECK_PARAMS
User=@PCP_USER@
