# src/pmie/pmie_daily.service.in
[Unit]
Description=Process PMIE logs
Documentation=man:pmie_daily(1)
ConditionPathExists=!@CRONTAB_PATH@

[Service]
Type=oneshot
Environment="PMIE_DAILY_PARAMS=-X xz -x 3"
EnvironmentFile=-@PCP_SYSCONFIG_DIR@/pmie_timers
ExecStart=@PCP_BINADM_DIR@/pmie_daily $PMIE_DAILY_PARAMS
User=@PCP_USER@

[Install]
RequiredBy=pmie.service
