# src/pmie/pmie_check.path.in
[Path]
PathChanged=@PCP_PMIECONTROL_PATH@
PathChanged=@PCP_PMIECONTROL_PATH@.d
