# src/pmie/pmie_daily.timer
[Unit]
Description=Daily processing of PMIE logs

[Timer]
OnCalendar=*-*-* 00:08:00

[Install]
WantedBy=timers.target
RequiredBy=pmie.service
