# src/pmie/pmie.service.in
[Unit]
Description=Performance Metrics Inference Engine
Documentation=man:pmie(1)
After=network-online.target pmcd.service
Before=pmie_farm.service pmie_check.timer pmie_daily.timer
BindsTo=pmie_farm.service pmie_check.timer pmie_daily.timer
Wants=pmcd.service

[Service]
Type=notify
NotifyAccess=all
# The pmie service manages the primary pmie (only)
Environment="PMIE_CHECK_PARAMS=--only-primary"
EnvironmentFile=-@PCP_SYSCONFIG_DIR@/pmie
TimeoutStartSec=120
TimeoutStopSec=120
Restart=always
ExecStart=@PCP_RC_DIR@/pmie start-systemd
ExecStop=@PCP_RC_DIR@/pmie stop-systemd
PIDFile=@PCP_RUN_DIR@/pmie.pid
Group=@PCP_GROUP@
User=@PCP_USER@

[Install]
WantedBy=multi-user.target
