# build/ci/playbooks/roles/azure-pipelines-agent/templates/azure-pipelines-agent.service.j2
[Unit]
Description=Azure Pipelines Agent
After=network.target

[Service]
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=5min
OOMScoreAdjust=-999
User=pipelines-agent
WorkingDirectory=/opt/azure-pipelines-agent
ExecStartPre=/opt/azure-pipelines-agent/config.sh --unattended
ExecStartPre=/bin/sh -c 'echo Agent.Image={{ image_name }} >> /opt/azure-pipelines-agent/.env'
ExecStartPre=/bin/sh -c 'curl -s -H Metadata:true "http://169.254.169.254/metadata/instance/compute/tags?api-version=2018-10-01&format=text" | \
                         tr ";" "\n" | grep Build. | tr ":" "=" >> /opt/azure-pipelines-agent/.env'
ExecStart=/opt/azure-pipelines-agent/bin/runsvc.sh
ExecStopPost=/opt/azure-pipelines-agent/config.sh remove --unattended

[Install]
WantedBy=multi-user.target
