# src/pmcd/pmcd.service.in
[Unit]
Description=Performance Metrics Collector Daemon
Documentation=man:pmcd(1)
After=network-online.target avahi-daemon.service
Before=zabbix-agent.service
 
[Service]
Type=notify
NotifyAccess=all
Restart=always
ExecStart=@PCP_RC_DIR@/pmcd start-systemd
ExecStop=@PCP_RC_DIR@/pmcd stop-systemd
PIDFile=@PCP_RUN_DIR@/pmcd.pid
 
[Install]
WantedBy=multi-user.target
