# src/pmlogger/pmlogger_check.service.in
[Unit]
Description=Check pmlogger instances are running
Documentation=man:pmlogger_check(1)
ConditionPathExists=!@CRONTAB_PATH@

[Service]
Type=oneshot
Environment="PMLOGGER_CHECK_PARAMS=-C"
EnvironmentFile=-@PCP_SYSCONFIG_DIR@/pmlogger_timers
ExecStart=@PCP_BINADM_DIR@/pmlogger_check $PMLOGGER_CHECK_PARAMS
User=@PCP_USER@

[Install]
RequiredBy=pmlogger.service
