# src/pmlogger/pmlogger_daily.timer
[Unit]
Description=Daily processing of archive logs

[Timer]
OnCalendar=*-*-* 00:10:00

[Install]
WantedBy=timers.target
RequiredBy=pmlogger.service
