# src/pmlogger/pmlogger_daily.service.in
[Unit]
Description=Process archive logs
Documentation=man:pmlogger(1)
ConditionPathExists=!@CRONTAB_PATH@

[Service]
Type=oneshot
ExecStart=@PCP_BINADM_DIR@/pmlogger_daily
User=@PCP_USER@

[Install]
RequiredBy=pmlogger.service
