# src/pmlogger/pmlogger_daily_poll.timer
[Unit]
Description=Hourly polling of log processing

[Timer]
OnCalendar=*-*-* *:30:00

[Install]
WantedBy=timers.target
RequiredBy=pmlogger.service
