# src/pmlogger/pmlogger_daily_poll.service.in
[Unit]
Description=Poll log processing
Documentation=man:pmlogger(1)

[Service]
Type=oneshot
ExecStart=@PCP_BINADM_DIR@/pmlogger_daily -p
User=@PCP_USER@
