# src/pmlogger/pmlogger_daily_report-poll.service.in
[Unit]
Description=Poll system activity report generation
Documentation=man:pmlogger_daily_report(1)
ConditionPathExists=!@CRONTAB_DAILY_REPORT_PATH@

[Service]
Type=oneshot
ExecStart=@PCP_BINADM_DIR@/pmlogger_daily_report -o @PCP_SA_DIR@ -p
User=@PCP_USER@

[Install]
WantedBy=pmlogger.service
RequiredBy=pmlogger_daily_report.service
