# src/pmlogger/pmlogger_check.service.in
[Unit]
Description=PCP archive loggers (non-primary and remote hosts)
Documentation=man:pmlogger_check(1)
ConditionPathExists=!@CRONTAB_PATH@

[Service]
Type=@SD_SERVICE_TYPE@
Restart=no
KillMode=process
TimeoutStartSec=4h
TimeoutStopSec=120
Environment="PMLOGGER_CHECK_PARAMS=--skip-primary"
EnvironmentFile=-@PCP_SYSCONFIG_DIR@/pmlogger_timers
ExecStart=@PCP_BINADM_DIR@/pmlogger_check --noerror $PMLOGGER_CHECK_PARAMS
WorkingDirectory=@PCP_VAR_DIR@
Group=@PCP_GROUP@
User=@PCP_USER@
