# src/pmlogger/pmlogger_check.path.in
[Path]
PathChanged=@PCP_PMLOGGERCONTROL_PATH@
PathChanged=@PCP_PMLOGGERCONTROL_PATH@.d
