# src/pmlogger/pmlogger_daily_check.service.in
[Unit]
Description=Check pmlogger instances are running
Documentation=man:pmlogger(1)
ConditionPathExists=!@CRONTAB_PATH@

[Service]
Type=oneshot
ExecStart=@PCP_BINADM_DIR@/pmlogger_daily -C
User=@PCP_USER@

[Install]
RequiredBy=pmlogger.service
