# src/pmlogger/pmlogger_daily.service.in
[Unit]
Description=Daily processing of local PCP archives
Documentation=man:pmlogger_daily(1)
ConditionPathExists=!@CRONTAB_PATH@

[Service]
Type=@SD_SERVICE_TYPE@
Restart=no
TimeoutStartSec=1h
Environment="PMLOGGER_DAILY_PARAMS=-E"
EnvironmentFile=-@PCP_SYSCONFIG_DIR@/pmlogger_timers
ExecStart=@PCP_BINADM_DIR@/pmlogger_daily --noerror -P $PMLOGGER_DAILY_PARAMS
WorkingDirectory=@PCP_VAR_DIR@
Group=@PCP_GROUP@
User=@PCP_USER@
