# src/pmlogger/pmlogger_daily_report.service.in
[Unit]
Description=Generate system activity reports
Documentation=man:pmlogger_daily_report(1)
ConditionPathExists=!@CRONTAB_DAILY_REPORT_PATH@

[Service]
Type=oneshot
Environment="PMLOGGER_DAILY_REPORT_PARAMS=-o @PCP_SA_DIR@"
EnvironmentFile=-@PCP_SYSCONFIG_DIR@/pmlogger_timers
ExecStart=@PCP_BINADM_DIR@/pmlogger_daily_report $PMLOGGER_DAILY_REPORT_PARAMS
User=@PCP_USER@

[Install]
WantedBy=pmlogger.service
