# src/pmlogger/pmlogger_farm_check.service.in
[Unit]
Description=Check and migrate non-primary pmlogger farm instances
Documentation=man:pmlogctl(1)
ConditionPathExists=!@CRONTAB_PATH@

[Service]
Type=exec
Restart=no
TimeoutStartSec=4h
TimeoutStopSec=120
ExecStart=@PCP_BIN_DIR@/pmlogctl -m check
WorkingDirectory=@PCP_VAR_DIR@
# root so pmlogctl can migrate pmlogger processes to the pmlogger_farm service
Group=root
User=root
