# src/pmlogger/pmlogger_check.service.in
[Unit]
Description=Check pmlogger instances are running
Documentation=man:pmlogger_check(1)
ConditionPathExists=!@CRONTAB_PATH@

[Service]
Type=oneshot
ExecStart=@PCP_BINADM_DIR@/pmlogger_check -C
User=@PCP_USER@

[Install]
RequiredBy=pmlogger.service
