# src/pmlogger/pmlogger_farm_check.service.in
[Unit]
Description=Check PCP archive loggers
Documentation=man:pmlogctl(1)
ConditionPathExists=!@CRONTAB_PATH@

[Service]
Type=@SD_SERVICE_TYPE@
Restart=no
TimeoutStartSec=4h
TimeoutStopSec=120
# needs to run as root:root not pcp:pcp, pmlogctl will handle the
# switch to pcp:pcp if any pmlogger needs to be launched, so NO
# User= nor Group= here
# -m is for migrate that we don't really need anymore
# check will check and if necessary restart all pmlogger instances
# known to pmlogctl(1)
ExecStart=@PCP_BIN_DIR@/pmlogctl -m check
WorkingDirectory=@PCP_VAR_DIR@
