# src/pmlogger/pmlogger_daily-poll.service.in
[Unit]
Description=Poll log processing
Documentation=man:pmlogger_daily(1)
ConditionPathExists=!@CRONTAB_PATH@

[Service]
Type=oneshot
ExecStart=@PCP_BINADM_DIR@/pmlogger_daily -p
User=@PCP_USER@

[Install]
RequiredBy=pmlogger.service
