# src/pmlogger/pmlogger_farm_check.service.in
[Unit]
Description=Check and migrate non-primary pmlogger farm instances
Documentation=man:pmlogctl(1)
ConditionPathExists=!@CRONTAB_PATH@

[Service]
Type=@SD_SERVICE_TYPE@
Restart=no
TimeoutStartSec=4h
TimeoutStopSec=120
# needs to run as root:root not pcp:pcp, pmlogctl will handle the
# switch to pcp:pcp if any pmlogger needs to be launched, so NO
# User= nor Group= here
ExecStart=@PCP_BIN_DIR@/pmlogctl -m check
WorkingDirectory=@PCP_VAR_DIR@
