# src/pmlogger/pmlogger_daily_report_poll.service.in
[Unit]
Description=Poll system activity report generation
Documentation=man:pmlogger(1)

[Service]
Type=oneshot
ExecStart=@PCP_BINADM_DIR@/pmlogger_daily_report -o @PCP_SA_DIR@ -p
User=@PCP_USER@
