# src/pmlogger/pmlogger_daily.service.in
[Unit]
Description=Process archive logs
Documentation=man:pmlogger_daily(1)
ConditionPathExists=!@CRONTAB_PATH@

[Service]
Type=oneshot
EnvironmentFile=-@PCP_SYSCONFIG_DIR@/pmlogger_timers
ExecStart=@PCP_BINADM_DIR@/pmlogger_daily $PMLOGGER_DAILY_PARAMS
User=@PCP_USER@

[Install]
RequiredBy=pmlogger.service
