# src/pmwebapi/pmwebd.service.in
[Unit]
Description=Performance Metrics Web Daemon
Documentation=man:pmwebd(8)
Wants=avahi-daemon.service
After=network.target avahi-daemon.service pmcd.service
 
[Service]
Type=oneshot
ExecStart=@path@/pmwebd start
ExecStop=@path@/pmwebd stop
RemainAfterExit=yes 
 
[Install]
WantedBy=multi-user.target
