# debian/pcp-webapi.postrm
#!/bin/sh -e
#DEBHELPER#

if [ "$1" = purge ]; then
    update-rc.d pmwebd remove >/dev/null
fi
