# debian/pcp.postrm
#!/bin/sh
#DEBHELPER#
if [ "$1" = purge ]; then
    update-rc.d pmproxy remove >/dev/null
    update-rc.d pmie remove >/dev/null
    update-rc.d pcp remove >/dev/null
fi
