# build/mac/io.pcp.pmcd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>io.pcp.pmcd</string>
        <key>ProgramArguments</key>
        <array>
            <string>/Library/StartupItems/pcp/pmcd</string>
            <string>/start</string>
        </array>
        <key>Disabled</key>
        <false/>
        <key>KeepAlive</key>
        <true/>
        <key>ThrottleInterval</key>
        <integer>30</integer>
        <key>RunAtLoad</key>
        <true/>
    </dict>
</plist>
