# build/mac/io.pcp.pmie_farm_check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>io.pcp.pmie_farm_check</string>
        <key>ProgramArguments</key>
        <array>
            <string>/usr/local/bin/pmiectl</string>
            <string>-m</string>
            <string>check</string>
        </array>
        <key>Disabled</key>
        <false/>
        <!-- Runs as root - pmiectl handles privilege switching to pcp:pcp -->
        <key>UserName</key>
        <string>root</string>
        <key>GroupName</key>
        <string>wheel</string>
        <key>StartCalendarInterval</key>
        <array>
            <dict>
                <key>Minute</key>
                <integer>29</integer>
            </dict>
            <dict>
                <key>Minute</key>
                <integer>59</integer>
            </dict>
        </array>
        <key>StandardErrorPath</key>
        <string>/var/log/pcp/pmie/pmie_farm_check.stderr</string>
        <key>StandardOutPath</key>
        <string>/var/log/pcp/pmie/pmie_farm_check.stdout</string>
        <key>WorkingDirectory</key>
        <string>/private/var/lib/pcp</string>
    </dict>
</plist>
