# build/mac/io.pcp.pmie.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
    pmie Health Check Service

    This is NOT a long-running daemon. It's a periodic health-check that:
    1. Runs at boot (RunAtLoad) and every 10 minutes (StartInterval)
    2. Spawns pmie processes if they're not running
    3. Verifies they're healthy, then exits

    Linux equivalent: pmie_check.timer + pmie_check.service (systemd)

    To reduce crash recovery time, decrease StartInterval (seconds).
    Minimum recommended: 60 (1 minute). Default: 600 (10 minutes).
-->
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>io.pcp.pmie</string>
        <key>ProgramArguments</key>
        <array>
            <string>/etc/init.d/pmie</string>
            <string>start</string>
        </array>
        <key>Disabled</key>
        <false/>
        <key>UserName</key>
        <string>root</string>
        <key>UserGroup</key>
        <string>wheel</string>
        <key>KeepAlive</key>
        <false/>
        <key>ThrottleInterval</key>
        <integer>30</integer>
        <key>RunAtLoad</key>
        <true/>
        <key>StartInterval</key>
        <integer>600</integer>
        <key>StandardErrorPath</key>
        <string>/var/log/pcp/pmie/plist.stderr</string>
        <key>StandardOutPath</key>
        <string>/var/log/pcp/pmie/plist.stdout</string>
    </dict>
</plist>
