# build/mac/io.pcp.pmlogger_daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>io.pcp.pmlogger_daily</string>
        <key>ProgramArguments</key>
        <array>
            <string>/usr/local/libexec/pcp/bin/pmlogger_daily</string>
            <string>--noerror</string>
            <string>-P</string>
        </array>
        <key>Disabled</key>
        <false/>
        <key>UserName</key>
        <string>pcp</string>
        <key>GroupName</key>
        <string>pcp</string>
        <key>StartCalendarInterval</key>
        <dict>
            <key>Hour</key>
            <integer>0</integer>
            <key>Minute</key>
            <integer>10</integer>
        </dict>
        <key>StandardErrorPath</key>
        <string>/var/log/pcp/pmlogger/pmlogger_daily.stderr</string>
        <key>StandardOutPath</key>
        <string>/var/log/pcp/pmlogger/pmlogger_daily.stdout</string>
        <key>WorkingDirectory</key>
        <string>/private/var/lib/pcp</string>
    </dict>
</plist>
