# src/pmfind/pmfind.path.in
[Path]
PathChanged=@PCP_SYSCONF_DIR@/pmfind

[Install]
RequiredBy=pmfind.service
