# build/containers/archive-analysis/root/etc/systemd/system/archive-import.service
[Unit]
Description=PCP archives import
Requires=valkey.service
After=valkey.service

[Service]
Type=exec
Restart=on-failure
StandardOutput=journal+console
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/archive-import.py
SyslogIdentifier=import

[Install]
WantedBy=multi-user.target
