# src/pcp/atop/pcp-atop-daily.service.in
[Unit]
Description=Daily pcp-atop archive rotation and maintenance
Documentation=man:pcp-atop(1)

[Service]
Type=oneshot
EnvironmentFile=-@PCP_SYSCONF_DIR@/atop/atop.conf
Environment="LOGPATH=@PCP_LOG_DIR@/atop"
Environment="LOGGENERATIONS=28"
ExecStart=systemctl stop pcp-atop.service
ExecStart=@PCP_BINADM_DIR@/pcp-atop-daily
ExecStart=systemctl start pcp-atop.service
