# src/pcp/atop/pcp-atop-rotate.service.in
[Unit]
Description=Restart pcp-atop to rotate archive logs
Documentation=man:pcp-atop(1)

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart pcp-atop.service
